# config/caddy.service
[Unit]
Description=Caddy

[Service]
ExecStart=caddy --conf /etc/caddy/Caddyfile

[Install]
WantedBy=multi-user.target

